# root/etc/s6-overlay/s6-rc.d/init-qbittorrent-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# make our folder
mkdir -p /config/qBittorrent

CONFIG_FILE="/config/qBittorrent/qBittorrent.conf"

# copy default config
if [[ ! -f "$CONFIG_FILE" ]]; then
    cp /defaults/qBittorrent.conf "$CONFIG_FILE"
fi

if [[ -z ${LSIO_NON_ROOT_USER} ]]; then
    # chown download directory if currently not set to abc
    if grep -qe ' /downloads ' /proc/mounts; then
        lsiown abc:abc /downloads
    fi

    # permissions
    lsiown -R abc:abc \
        /config
fi

# function to generate PBKDF2 hash exactly like qBittorrent source code
generate_pbkdf2_hash() {
    local password="$1"
    echo "$password" | python3 -c "
import hashlib
import os
import base64
import sys

# Read password from stdin
password = sys.stdin.read().rstrip('\n')

# Generate 16 bytes of random salt
salt = os.urandom(16)

# Generate PBKDF2 hash with exactly qBittorrent's parameters
hash_bytes = hashlib.pbkdf2_hmac('sha512', password.encode('utf-8'), salt, 100000, 64)

# Format exactly like qBittorrent: saltView.toBase64() + ':' + outBufView.toBase64()
salt_b64 = base64.b64encode(salt).decode('ascii')
hash_b64 = base64.b64encode(hash_bytes).decode('ascii')

print(f'{salt_b64}:{hash_b64}')
"
}

# ensure [Preferences] section exists if WebUI credentials are being configured
if [[ -n ${WEBUI_USER} ]] || [[ -n ${WEBUI_PASS} ]]; then
    if ! grep -q "^\s*\[\s*Preferences\s*\]" "$CONFIG_FILE"; then
        echo "" >> "$CONFIG_FILE"
        echo "[Preferences]" >> "$CONFIG_FILE"
    fi
fi

# configure WebUI password
if [[ -n ${WEBUI_PASS} ]]; then
    echo "setting WebUI password using WEBUI_PASS"
    PBKDF2_HASH=$(generate_pbkdf2_hash "$WEBUI_PASS")
    sed -i '/^\s*WebUI\\Password_PBKDF2\s*=/d' "$CONFIG_FILE"
    sed -i '/^\s*\[\s*Preferences\s*\]/a WebUI\\Password_PBKDF2="@ByteArray('"$PBKDF2_HASH"')"' "$CONFIG_FILE"
fi

# configure WebUI username
if [[ -n ${WEBUI_USER} ]]; then
    echo "setting WebUI username using WEBUI_USER"
    sed -i '/^\s*WebUI\\Username\s*=/d' "$CONFIG_FILE"
    sed -i '/^\s*\[\s*Preferences\s*\]/a WebUI\\Username="'"$WEBUI_USER"'"' "$CONFIG_FILE"
fi
